# deploy/residential-proxy.service
[Unit]
Description=Residential SOCKS5 proxy forwarder (gost to upstream residential IP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.openclaw/residential-proxy.env
ExecStart=/usr/local/bin/gost -L socks5://127.0.0.1:1080 -F "socks5://${IPROYAL_USER}:${IPROYAL_PASS}@${IPROYAL_HOST}:${IPROYAL_PORT}"
Restart=always
RestartSec=10
StandardOutput=append:%h/.openclaw-oauth/residential-proxy.log
StandardError=append:%h/.openclaw-oauth/residential-proxy.log
OnFailure=residential-proxy-alert.service

[Install]
WantedBy=default.target
